# infrastructure/systemd/dis-tts-engine.service
[Unit]
Description=DIS Self-hosted Speech Engine
Wants=network-online.target
After=network-online.target
Before=dis-speech.service

[Service]
Type=exec
User=dis
Group=dis
WorkingDirectory=/opt/dis/speech-engine
RuntimeDirectory=dis-tts
RuntimeDirectoryMode=0750
UMask=0027
ExecStartPre=/usr/bin/test -x /opt/dis-data/tts/runtime/bin/python
ExecStartPre=/usr/bin/test -d /opt/dis-data/tts/models
ExecStartPre=/usr/bin/test -d /opt/dis-data/tts/staging
ExecStartPre=/usr/bin/test -d /opt/dis-data/tts/state
ExecStartPre=/usr/bin/setfacl -m u:www-data:--x,d:u:www-data:rw- /run/dis-tts
ExecStart=/opt/dis-data/tts/runtime/bin/python -m dis_tts_engine
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
KillMode=control-group
TimeoutStopSec=60s
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=DIS_TTS_SOCKET_PATH=/run/dis-tts/engine.sock
Environment=DIS_TTS_MODELS_ROOT=/opt/dis-data/tts/models
Environment=DIS_TTS_STAGING_ROOT=/opt/dis-data/tts/staging
Environment=DIS_TTS_STATE_ROOT=/opt/dis-data/tts/state
Environment=DIS_TTS_SYNTHESIS_DEADLINE_SECONDS=14300
Environment=DIS_TTS_TORCH_THREADS=16
Environment=OMP_NUM_THREADS=16
Environment=OMP_THREAD_LIMIT=16
Environment=MKL_NUM_THREADS=16
Environment=OPENBLAS_NUM_THREADS=16
Environment=NUMEXPR_NUM_THREADS=16
Environment=TOKENIZERS_PARALLELISM=false
Environment=HF_HOME=/opt/dis-data/tts/models/.hf-cache
Environment=TORCH_HOME=/opt/dis-data/tts/models/.torch-cache
Environment=XDG_CACHE_HOME=/opt/dis-data/tts/models/.runtime-cache
Environment=TORCHINDUCTOR_CACHE_DIR=/opt/dis-data/tts/models/.torchinductor-cache
Environment=HF_HUB_DISABLE_TELEMETRY=1
Environment=HF_HUB_DISABLE_IMPLICIT_TOKEN=1
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectProc=invisible
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
RestrictRealtime=true
ReadWritePaths=/opt/dis-data/tts/models /opt/dis-data/tts/staging /opt/dis-data/tts/state
ReadOnlyPaths=/opt/dis/speech-engine /opt/dis-data/tts/runtime /opt/dis-data/tts/python
InaccessiblePaths=-/opt/dis-data/backup -/opt/dis-data/secrets -/opt/dis-data/storage -/opt/dis-data/webapp/backend/storage
MemoryHigh=80%
MemoryMax=90%
CPUQuota=1600%
TasksMax=512
LimitNOFILE=65536
Nice=10
OOMPolicy=stop

[Install]
WantedBy=multi-user.target
