# infrastructure/systemd/dis-osrm.service
[Unit]
Description=DIS local OSRM routing engine
Documentation=https://github.com/Project-OSRM/osrm-backend
After=network.target
ConditionPathExists=@DIS_DATA_PATH@/osrm/current/routing.osrm.partition

[Service]
Type=simple
User=dis-osrm
Group=dis-osrm
Environment=DIS_DATA_PATH=@DIS_DATA_PATH@
Environment=APP_ROOT=@APP_ROOT@
ExecStart=/usr/bin/bash @OSRM_RUNTIME_SCRIPT@ serve
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=30
LimitNOFILE=65536
TasksMax=128
MemoryHigh=6G
MemoryMax=8G
CPUQuota=300%
OOMPolicy=stop
UMask=0027

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictAddressFamilies=AF_INET AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
IPAddressDeny=any
IPAddressAllow=localhost

[Install]
WantedBy=multi-user.target
