# infrastructure/systemd/dis-knmi-realtime.service
[Unit]
Description=DIS KNMI Realtime Precipitation Import Worker
After=network-online.target redis-server.service postgresql.service
Wants=network-online.target

[Service]
Type=exec
User=dis
Group=dis
WorkingDirectory=/opt/dis/webapp/backend
ExecStartPre=/usr/bin/test -x /usr/bin/h5dump
ExecStartPre=/usr/bin/test -r /opt/dis/webapp/backend/artisan
ExecStart=/usr/bin/php -d zend.exception_ignore_args=1 artisan queue:work knmi_realtime --queue=knmi-realtime --sleep=2 --tries=1 --timeout=1200 --max-jobs=250 --max-time=86400
Restart=always
RestartSec=5
KillSignal=SIGTERM
KillMode=control-group
TimeoutStopSec=30s
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
RestrictRealtime=true
MemoryHigh=384M
MemoryMax=512M
CPUQuota=100%
TasksMax=32
Nice=5
OOMPolicy=stop
ReadWritePaths=/opt/dis-data /opt/dis/webapp/backend/storage /opt/dis/storage

[Install]
WantedBy=multi-user.target
